# deploy/dbrain-doctor.service
[Unit]
Description=d-brain daily self-diagnostic
OnFailure=dbrain-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/projects/dbrain
ExecStart=%h/.local/bin/uv run python -m d_brain.services.doctor
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_CONFIG_DIR=%h/.claude
