# deploy/dbrain-doctor.timer
[Unit]
Description=d-brain daily self-diagnostic (08:00 local)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
